# packageFiles/postinst
#!/bin/bash

cd /apps
npm install
chmod +x start.sh
systemctl enable sd-device.service
systemctl start sd-device.service