# packageFiles/postinst
#!/bin/bash

apt-get update
apt-get install npm -y
npm install n -g
n stable
cd /apps
npm install n -g
chmod +x start.sh
systemctl enable sd-device.service
systemctl start sd-device.service