# packageFiles/postinst
#!/bin/bash

npm install n -g
n stable
cd /apps
npm install
chmod +x start.sh
systemctl enable sd-device.service
systemctl start sd-device.service